# supervisord_conf/conf.d/http.ini
[program:httpd]
command=/usr/sbin/httpd -DFOREGROUND
autostart=true
autorestart=true
startsecs=5
startretries=3
user=www-data
stopsignal=KILL
stopwaitsecs=0
killasgroup=true
stopasgroup=true
