# supervisord_conf/conf.d/php-fpm.ini
[program:php-fpm]
command=/usr/local/sbin/php-fpm -g /run/php-fpm.pid -F
autostart=true
autorestart=unexpected
startsecs=5
startretries=3
user=www-data
stopsignal=KILL
stopwaitsecs=0
killasgroup=true
stopasgroup=true
stderr_logfile=/var/log/supervisord.log
